# daemon/systemd/cronplus.service
[Unit]
Description=cronplus - Advanced Cron Task Manager
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/cronplus/src/main.py --conf /opt/cronplus/tasks.conf --logs /opt/cronplus/logs/logs.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
User=root

# 安全加固（保留与 root 运行兼容的选项）
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/opt/cronplus /run
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
